# status-email@.service
[Unit]
Description=Sendet eine E-Mail bei Service-Fehlern

[Service]
Type=oneshot
ExecStart=/usr/local/bin/systemd-email.sh %i
